# sdcard/etc/rc.local
\ this file will be included on startup
cr .( include 0:/etc/rc.local ) cr

