# sdcard/etc/rc.local
\ this file will be included on startup
cr .( include 0:/etc/rc.local ) cr

\ watchdog? [IF] .( Watchdog has bitten, # bites: ) watchdog# . cr [THEN]
\ assert? [IF] assert# . .( Assertations since cold startup, last one: ) assert@ drop .assert [THEN]
